# BlogProject/crontab
* * * * * /usr/local/bin/python /usr/src/BlogProject/blog-cron.py >> /var/log/cron.log 2>&1
# Ensure there is a newline at the end of the file





